# linux/config/systemd/user/isync.timer
[Timer]
OnActiveSec=0s
OnUnitInactiveSec=300s
Unit=isync.service

[Install]
WantedBy=timers.target
