# linux/config/systemd/user/offlineimap.timer
[Timer]
OnActiveSec=0s
OnUnitInactiveSec=300s
Unit=offlineimap.service

[Install]
WantedBy=timers.target
